# rclone-auto-sync.timer
[Unit]
Description=Run rclone-auto-sync every minute

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target